# puppet/puppet-1604/modules/system-nas/files/etc/cron.d/backup-rusa

15 3 * * * root /usr/local/bin/backup-linode-rusa.sh
30 3 * * * root /usr/local/bin/backup-linode-rusa-www.sh
