# puppet/puppet-1604/modules/system-nas/files/etc/cron.d/backup-linode

0 2 * * * root /usr/local/bin/backup-linode.sh
