# puppet/puppet-1410/modules/system-nas/files/etc/cron.d/backup-linode

0 2 * * * root /usr/local/bin/backup-linode.sh
15 2 * * * root /usr/local/bin/backup-linode2.sh
