# puppet/puppet-1404/modules/system-nas/files/etc/cron.d/backup-folders

0 4 * * * root /usr/local/bin/backup.sh cds
10 4 * * * root /usr/local/bin/backup.sh family
20 4 * * * root /usr/local/bin/backup.sh radio
30 4 * * * root /usr/local/bin/backup.sh DPBARCHIVE
40 4 * * * root /usr/local/bin/backup-mirror.sh
