# puppet/puppet-1404/modules/system-nas/files/etc/cron.d/backup-linode

15 2 * * * root /usr/local/bin/backup-linode2.sh
